# dot_config/systemd/user/batteries.service
[Unit]
Description=Notify when mouse batteries are low

[Service]
Type=simple
ExecStart=/home/emil/.local/scripts/batteries

[Install]
WantedBy=default.target
