# dot_config/systemd/user/download_calendars.timer
[Unit]
Description=Download shared calendars

[Timer]
OnCalendar=00/2:00
Persistent=True

[Install]
WantedBy=timers.target
