# dot_config/systemd/user/home-emil-docs-blog.automount
[Unit]
Description=Automount SSHFS blog

[Automount]
Where=/home/emil/docs/blog

[Install]
WantedBy=default.target
