# dot_config/systemd/user/morning-flag.service
[Unit]
Description=Reset morning flag daily

[Service]
Type=oneshot
ExecStart=/bin/rm -f /tmp/morning

[Install]
WantedBy=default.target

