# dot_config/systemd/user/morning.timer
[Unit]
Description=Run morning routine on first boot of day

[Timer]
OnBootSec=0
OnCalendar=06:00
Persistent=true

[Install]
WantedBy=timers.target
