# dot_config/systemd/user/newsraft.service
[Unit]
Description=Update RSS feeds

[Service]
Type=oneshot
ExecStart=newsraft -e reload-all

[Install]
WantedBy=default.target
