# dot_config/systemd/user/noop.service
[Unit]
Description=No-Op Service

[Service]
Type=oneshot
ExecStart=/bin/true
