# dot_config/systemd/user/notify-failed@.service
[Unit]
Description=Send notifications when other services fail

[Service]
Type=simple
ExecStart=notify-send -u normal "%I failed"
