# dot_config/systemd/user/udiskie.service
[Unit]
Description=Automount new drives

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
