# dot_config/systemd/user/vdirsyncer@.service
[Unit]
Description=Sync %i calendar hourly

[Service]
Type=simple
ExecStart=vdirsyncer sync personal_calendar/%i

[Install]
WantedBy=

