# dot_config/systemd/user/vdirsyncer_gcal@.service
[Unit]
Description=Sync Google %i calendar hourly

[Service]
Type=simple
ExecStart=vdirsyncer sync gcal/%i

[Install]
WantedBy=

