# dot_config/vim/skeletons/skeleton.timer
[Unit]
Description=

[Timer]
OnCalendar=
Persistent=True

[Install]
WantedBy=timers.target
