# dot_config/vim/skeletons/skeleton.service
[Unit]
Description=

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=default.target
